# imageroot/systemd/user/redis_main-app.service
#
# Copyright (C) 2022 Nethesis S.r.l.
# SPDX-License-Identifier: GPL-3.0-or-later
#

[Unit]
Description=Podman  redis_main-app.service
BindsTo=kobotoolbox.service
After=kobotoolbox.service

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
EnvironmentFile=%S/state/environment
EnvironmentFile=-%S/state/smarthost.env
WorkingDirectory=%S/state
Restart=always
TimeoutStopSec=70
ExecStartPre=/bin/mkdir -p tmp
ExecStartPre=/bin/rm -f %t/redis_main-app.pid %t/redis_main-app.ctr-id
ExecStartPre=-runagent discover-smarthost
ExecStart=/usr/bin/podman run --conmon-pidfile %t/redis_main-app.pid  --cidfile %t/redis_main-app.ctr-id --cgroups=no-conmon  --pod-id-file %t/kobotoolbox.pod-id --replace -d --name  redis_main-app     --env redis_main.env  ${REDIS_IMAGE}  
ExecStop=/usr/bin/podman stop --ignore --cidfile %t/redis_main-app.ctr-id -t 10
ExecReload=/usr/bin/podman kill -s HUP redis_main-app
SyslogIdentifier=%u
ExecStopPost=/usr/bin/podman rm --ignore -f --cidfile %t/redis_main-app.ctr-id
PIDFile=%t/redis_main-app.pid
Type=forking

[Install]
WantedBy=default.target
